# deployments/systemd/poly-oracle.service
[Unit]
Description=Poly Oracle - Event Probability Monitor
Documentation=https://github.com/poly-oracle
After=network.target

[Service]
Type=simple
User=poly-oracle
Group=poly-oracle
WorkingDirectory=/opt/poly-oracle
ExecStart=/opt/poly-oracle/poly-oracle --config /opt/poly-oracle/configs/config.yaml
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=poly-oracle

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
